# bin/lumi-lights.service
[Unit]
Description=Lumi Endurance Run

[Service]
ExecStart=/usr/bin/node /usr/local/lumi/src/entry/lights.js
Restart=always
RestartSec=10

# Output to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=lumi-lights-service
Environment=

User=root

[Install]
WantedBy=multi-user.target