# node/bin/lumi.service
[Unit]
Description=Lumi Zone Game
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/usr/local/lumi/node/.env
ExecStart=/usr/local/bin/node /usr/local/lumi/node/src/index.js
Restart=always
RestartSec=10

# Output to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=lumi-service
Environment=

User=root

[Install]
WantedBy=multi-user.target